# vendor/spoken_source/launchd/com.creativemachines.ostler.spoken-bundle.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!--
        Ostler spoken bundle LaunchAgent (HR015 conversation-memory feed).

        This is the four-artefact leg for SPOKEN conversations: meetings,
        calls, AND voice notes. It does NOT record or transcribe; CM042
        (RemoteCapture) owns capture + Whisper transcription and writes
        finished markdown transcripts under
        ~/Documents/Ostler/Transcripts/YYYY/MM/. This agent reads those
        finished transcripts, normalises each session into a cleaned
        transcript + CM048 metadata (channel="spoken"), and hands each
        new session to CM048's pwg-convo processor, which emits the four
        artefacts under
        ~/Documents/Ostler/Conversations/<date>/<slug>-<short-id>/.

        It sits alongside the email-bundle and imessage-bundle agents
        and shares their wiring shape. Distinct label, wrapper, state
        file, and output, so the three feeds never collide.

        Pieces:
          1. spoken_source/ is the reader + renderer + dispatcher.
          2. bin/spoken-bundle-tick.sh is the wrapper this agent runs;
             it drives the module against the RemoteCapture tree.
          3. CM051 install.sh copies the wrapper, renders this plist,
             and calls launchctl bootstrap.

        The installer renders the placeholders at install time:
          OSTLER_BIN         -> $OSTLER_DIR/bin
          OSTLER_HOME        -> $HOME (working directory)
          OSTLER_LOGS        -> $LOGS_DIR (~/.ostler/logs by default)
          OSTLER_PYTHON_PATH -> the spoken-source venv python3
          PWG_CONVO_CMD      -> absolute pwg-convo invocation
                                (the CM048 venv python -m src.cli, or
                                the installed console script)

        British English throughout.
    -->
    <key>Label</key>
    <string>com.creativemachines.ostler.spoken-bundle</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>OSTLER_BIN/spoken-bundle-tick.sh</string>
    </array>

    <!--
        Fifteen minutes, matching the email-bundle + imessage-bundle
        cadence. RemoteCapture writes a transcript only after a
        recording finishes, so a quarter-hour keeps recent meetings and
        voice notes fresh in the assistant's memory without re-walking
        the transcript tree too aggressively. The watermark in
        ~/.ostler/workspace/spoken_source_state.json means each tick
        only dispatches a session it has not bundled before.
    -->
    <key>StartInterval</key>
    <integer>900</integer>

    <!--
        RunAtLoad fires once at login / post-install. The wrapper
        clamps the read window via --since-days so a fresh install
        does not try to bundle every historic recording in one tick.
        The wrapper also exits 0 cleanly when the transcripts directory
        is absent (RemoteCapture never run), so RunAtLoad never
        hard-fails the agent on a Mac without any recordings.
    -->
    <key>RunAtLoad</key>
    <true/>

    <key>WorkingDirectory</key>
    <string>OSTLER_HOME</string>

    <key>StandardOutPath</key>
    <string>OSTLER_LOGS/spoken-bundle.log</string>
    <key>StandardErrorPath</key>
    <string>OSTLER_LOGS/spoken-bundle.err</string>

    <!--
        OSTLER_PYTHON selects the spoken-source venv interpreter.
        PWG_CONVO_CMD tells the dispatcher how to invoke CM048's CLI
        (rendered to the CM048 venv at install time). Without these the
        tick falls through to a stock system python with no
        site-packages -- the same launch-blocker shape the email-bundle
        agent guards against.
    -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>OSTLER_PYTHON</key>
        <string>OSTLER_PYTHON_PATH</string>
        <key>PWG_CONVO_CMD</key>
        <string>PWG_CONVO_CMD_VALUE</string>
        <!--
            Rendered by the installer (_install_conversation_feed) so the
            wrapper ships verbatim: SOURCE_DIR is the package parent the
            wrapper cd's into; USER_DISPLAY_NAME labels the operator's own
            speech as "You".
        -->
        <key>OSTLER_SOURCE_DIR</key>
        <string>OSTLER_SOURCE_DIR_VALUE</string>
        <key>OSTLER_USER_DISPLAY_NAME</key>
        <string>OSTLER_USER_DISPLAY_NAME_VALUE</string>
        <!--
            OSTLER_USER_ID scopes the conversation to the operator. CM048's
            pwg-convo fails loud at settings construction without it (the
            productisation guard against an empty user_id silently
            propagating), so the bundle agent must carry it or every tick
            dies before processing a single chat. Rendered from the
            installer's USER_ID.
        -->
        <key>OSTLER_USER_ID</key>
        <string>OSTLER_USER_ID_VALUE</string>
        <!--
            T16 at-rest encryption: pwg-convo (CM048) fails closed
            without a DB key, so the feed agent must carry it.
            OSTLER_DB_KEY_FILE points at the installer's mode-0600
            key file; OSTLER_ALLOW_PLAINTEXT is 1 only for a dev/CI
            --allow-plaintext install (else 0, inert).
        -->
        <key>OSTLER_DB_KEY_FILE</key>
        <string>OSTLER_DB_KEY_FILE_VALUE</string>
        <key>OSTLER_ALLOW_PLAINTEXT</key>
        <string>OSTLER_ALLOW_PLAINTEXT_VALUE</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>Nice</key>
    <integer>5</integer>
</dict>
</plist>
