# vendor/email_ingest/launchd/com.creativemachines.ostler.email-ingest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!--
        Ostler email ingest LaunchAgent.

        Hourly, drains any new messages from Apple Mail's local store
        into a Gmail-format mbox and hands the mbox to CM046's email
        adapter, which threads, cleans, and writes CM048 conversation
        files. The result: emails the user got in the last hour are
        in the assistant's knowledge graph by the next tick.

        Pieces:
          1. ostler_fda/apple_mail_mbox.py is the emit-mbox primitive.
          2. email-ingest/bin/email-ingest-tick.sh is the wrapper
             this LaunchAgent runs; it drives 1 plus pwg-email-ingest.
          3. CM051 install.sh integration copies the wrapper, renders
             this plist, and calls launchctl bootstrap.

        The installer renders the placeholders below at install time:
          OSTLER_BIN  -> $OSTLER_DIR/bin
          OSTLER_HOME -> $HOME (used as working directory)
          OSTLER_LOGS -> $LOGS_DIR (~/.ostler/logs by default)

        British English throughout.
    -->
    <key>Label</key>
    <string>com.creativemachines.ostler.email-ingest</string>

    <!--
        Wrapper script does both halves of the chain (emit plus ingest)
        so the LaunchAgent stays a single-process unit. The script is
        installed to OSTLER_BIN/email-ingest-tick.sh by CM051's
        installer; placeholder is replaced at install time.
    -->
    <key>ProgramArguments</key>
    <array>
        <string>OSTLER_ASSISTANT_BINARY</string>
        <string>run-source</string>
        <string>email-ingest</string>
    </array>

    <!--
        One hour. Apple Mail syncs in the background as it always
        does, so any reasonable cadence works; an hour gives the
        delivery and sync layer time to settle without leaving recent
        mail invisible to the assistant for too long.
    -->
    <key>StartInterval</key>
    <integer>3600</integer>

    <!--
        RunAtLoad fires the script once as soon as the agent loads
        (login or post-install). On a fresh install the wrapper
        clamps the emit window via the backfill-days flag so we
        don't try to ingest the user's entire mail history in one
        tick.
    -->
    <key>RunAtLoad</key>
    <true/>

    <!--
        Sane working directory: the wrapper is absolute-path safe
        but $HOME-relative resolution still works from here.
    -->
    <key>WorkingDirectory</key>
    <string>OSTLER_HOME</string>

    <!--
        Stdout and stderr go to $LOGS_DIR. The installer creates the
        dir. Doctor surfaces these logs so the operator can see
        whether the last tick succeeded without spelunking launchctl.
    -->
    <key>StandardOutPath</key>
    <string>OSTLER_LOGS/email-ingest.log</string>
    <key>StandardErrorPath</key>
    <string>OSTLER_LOGS/email-ingest.err</string>

    <!--
        OSTLER_PYTHON tells the tick script which python interpreter
        to use. Without this, the tick falls through to PATH-`python3`
        which on a stock macOS is /Library/Developer/CommandLineTools
        /usr/bin/python3 -- the system python with no installed
        site-packages. The tick imports ostler_fda.apple_mail_mbox,
        which only exists inside the dedicated email-ingest venv
        created by CM051's install.sh. Failing to set this caused
        the LAUNCH-BLOCKER ModuleNotFoundError seen in retest #N
        (2026-05-23). The installer renders OSTLER_PYTHON_PATH ->
        ${OSTLER_DIR}/services/email-ingest/.venv/bin/python3 at
        install time.
    -->
    <!--
        PWG_EMAIL_INGEST is the ABSOLUTE path to CM021's
        pwg-email-ingest console script for the ingest leg. Under
        launchd the agent inherits a minimal PATH
        (/usr/bin:/bin:/usr/sbin:/sbin) that does NOT include the
        email-ingest venv's bin/, so a bare `pwg-email-ingest` lookup
        always missed and EVERY harvested message was silently dropped
        before reaching the graph (the "Emitted N message(s)" then
        "not on PATH" / exit-127 shape). The installer renders
        PWG_EMAIL_INGEST_PATH ->
        ${OSTLER_DIR}/services/email-ingest/.venv/bin/pwg-email-ingest
        at install time, mirroring OSTLER_PYTHON_PATH above.
    -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>OSTLER_PYTHON</key>
        <string>OSTLER_PYTHON_PATH</string>
        <key>PWG_EMAIL_INGEST</key>
        <string>PWG_EMAIL_INGEST_PATH</string>
    </dict>

    <!--
        Background priority. Email ingest is not latency-critical;
        we don't want it competing with foreground work for I/O on
        the user's daily-driver Mac.
    -->
    <key>ProcessType</key>
    <string>Background</string>

    <key>Nice</key>
    <integer>5</integer>
</dict>
</plist>
