# assistant-agent/launchd/com.creativemachines.ostler.assistant.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!--
        Ostler assistant LaunchAgent (daemon mode).

        Long-running autonomous runtime: gateway server, configured
        channels (iMessage and email at v0.1), heartbeat monitor,
        and the cron scheduler. The customer-facing binary
        ostler-assistant is the upstream zeroclaw renamed at tar
        time; daemon semantics are unchanged.

        The installer renders the placeholders below at install time:
          OSTLER_BIN              -> $OSTLER_DIR/OstlerAssistant.app/Contents/MacOS
          OSTLER_HOME             -> $HOME
          OSTLER_LOGS             -> $LOGS_DIR
          OSTLER_ASSISTANT_CONFIG -> $OSTLER_DIR/assistant-config

        OSTLER_BIN points at the .app bundle's MacOS dir (v0.4.3+
        shape) instead of the legacy ~/.ostler/bin/. The bundle
        wrapper exists so macOS TCC + Activity Monitor + the FDA
        grant dialog render the Ostler v4 oxblood squircle next
        to the daemon. Runtime semantics are unchanged: launchd
        loads Contents/MacOS/ostler-assistant directly, the
        bundle wrapper does not affect the daemon's process model.

        ZEROCLAW_WORKSPACE points the binary at Phase D's config
        directory. The config-dir resolver short-circuits to the
        workspace value when config.toml exists there, so this
        single env var binds runtime to the installer-emitted
        TOML without touching the binary CLI surface.

        British English throughout.
    -->
    <key>Label</key>
    <string>com.creativemachines.ostler.assistant</string>

    <!--
        daemon subcommand spins up the long-running runtime
        (gateway + channels + heartbeat + scheduler). See
        ostler-ai/ostler-assistant src/main.rs Commands::Daemon.

        OSTLER_BIN expands to the inner-bundle MacOS dir
        (v0.4.3+) so launchd invokes
        OstlerAssistant.app/Contents/MacOS/ostler-assistant.
        macOS TCC sees the parent .app bundle ID (read from
        Info.plist CFBundleIdentifier) and renders the daemon
        with the Ostler v4 icon in System Settings.
    -->
    <key>ProgramArguments</key>
    <array>
        <string>OSTLER_BIN/ostler-assistant</string>
        <string>daemon</string>
    </array>

    <!--
        Config resolution via env. ZEROCLAW_WORKSPACE is the
        documented top-level entry point (resolve_runtime_config_dirs
        in crates/zeroclaw-config/src/schema.rs). When config.toml
        exists at the workspace path, resolve_config_dir_for_workspace
        treats the workspace dir AS the config dir. Phase D's
        installer emits config.toml at $OSTLER_DIR/assistant-config,
        so the env var below is sufficient.
    -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>ZEROCLAW_WORKSPACE</key>
        <string>OSTLER_ASSISTANT_CONFIG</string>
        <!--
            OSTLER_IMESSAGE_SELF_HANDLES (#646): the customer's OWN
            iMessage handles (comma-separated phone + email), rendered by
            INSTALL_SNIPPET from the me-card identity captured during
            install. The assistant's self-echo loop guard
            (crates/zeroclaw-channels/src/imessage.rs is_self_handle)
            reads this env var and drops inbound that matches one of the
            user's own handles, so the daemon cannot reply to its own
            output echoing back as inbound (shared Apple ID / cross-account
            routing). Empty when no handle was captured; the guard then
            stays inactive and the daemon's content-based backstop applies.
        -->
        <key>OSTLER_IMESSAGE_SELF_HANDLES</key>
        <string>OSTLER_IMESSAGE_SELF_HANDLES_VALUE</string>
        <!--
            OSTLER_IMESSAGE_SELF_CHAT (self-chat reachability, 2026-06-17):
            lets the customer text the assistant from their OWN Apple ID
            (Note to Self) and get a reply, with no second account. With this
            on, the chat.db poll also reads the operator's own messages
            (is_from_me=1) but ONLY within a self-chat (a conversation whose
            single participant is the account owner), scoped by a verified
            predicate that excludes every normal outgoing message
            (crates/zeroclaw-channels/src/imessage.rs POLL_QUERY_SELF_CHAT).
            The daemon refuses to enable it unless a non-empty identity prefix
            is in force, because that prefix is the echo guard that stops the
            assistant re-ingesting its own self-chat replies. Fixed product
            constant (the core "text your assistant" promise), so hard-coded
            here rather than rendered by INSTALL_SNIPPET.
        -->
        <key>OSTLER_IMESSAGE_SELF_CHAT</key>
        <string>true</string>
        <!--
            OLLAMA_NUM_CTX (.153 cold-wipe walk, 2026-06-16): the daemon's
            agent system prompt is ~7,000 tokens. Ollama's default context
            window is 4096, so without this the prompt is silently truncated
            to 4095 tokens and qwen returns an EMPTY completion, which the
            gateway surfaces to the user as "500 Internal Server Error /
            {error: EOF}". The daemon reads OLLAMA_NUM_CTX (and
            OLLAMA_KEEP_ALIVE) from its env, but nothing was setting it, so
            every real chat failed on a fresh install. 32768 fits the full
            prompt with headroom (~6.8 GB VRAM for qwen3.5:9b, well within
            a customer Mac). Proven live on the .153 box: at ctx 4096 the
            agent prompt returns empty; at 32768 it returns a real reply.
            KEEP_ALIVE pins the model resident so it does not cold-reload
            between turns. These are fixed product constants (not
            customer-variable), so they are hard-coded here rather than
            rendered by INSTALL_SNIPPET.
        -->
        <key>OLLAMA_NUM_CTX</key>
        <string>32768</string>
        <key>OLLAMA_KEEP_ALIVE</key>
        <string>30m</string>
    </dict>

    <!--
        RunAtLoad: start the daemon at LaunchAgent load (login
        or post-install).

        KeepAlive (dict form): only restart on abnormal exit. A
        clean exit zero (e.g. SIGTERM during uninstall, or daemon
        deciding to stop on a fatal config error it has already
        surfaced) does not retry. SuccessfulExit=false matches
        the launchd convention for "respawn on crash, leave alone
        on intentional shutdown".
    -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!--
        Throttle restarts so a tight crash loop on a misconfigured
        channel does not burn battery / saturate the log rotator.
        10 seconds matches the launchd default but is set
        explicitly so the value is reviewable in the file.
    -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!--
        Working directory: HOME so any relative paths in the
        binary resolve against the user's home rather than /.
    -->
    <key>WorkingDirectory</key>
    <string>OSTLER_HOME</string>

    <!--
        Stdout and stderr go to $LOGS_DIR. Doctor surfaces these
        logs so the operator can see boot-time failures without
        spelunking launchctl print.
    -->
    <key>StandardOutPath</key>
    <string>OSTLER_LOGS/ostler-assistant.log</string>
    <key>StandardErrorPath</key>
    <string>OSTLER_LOGS/ostler-assistant.err</string>

    <!--
        Interactive process type because the daemon serves
        websocket gateway traffic the user is interactively
        producing (iMessage replies, dashboard clicks). Adaptive
        is the closest match for "responsive but not realtime";
        but Background is what the existing wiki-recompile uses
        and is the conservative pick for a v0.1 always-on agent
        that is mostly idle. Open question for the operator:
        flip to Adaptive if iMessage latency feels sluggish.
    -->
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
