# vendor/imessage_source/launchd/com.creativemachines.ostler.imessage-bundle.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!--
        Ostler iMessage bundle LaunchAgent (CM040 product path).

        Mirrors the email-ingest tick. On an interval, reads the Hub
        Mac's iMessage chat.db (Full Disk Access granted at install),
        threads recent messages into conversation sessions, and hands
        each new session to CM048's pwg-convo processor, which emits
        the four artefacts under
        ~/Documents/Ostler/Conversations/<date>/<slug>-<short-id>/.

        Pieces:
          1. services/imessage_source/ is the threader + dispatcher.
          2. bin/imessage-bundle-tick.sh is the wrapper this agent
             runs; it drives the module against chat.db.
          3. CM051 install.sh copies the wrapper, renders this plist,
             and calls launchctl bootstrap.

        The installer renders the placeholders at install time:
          OSTLER_BIN         -> $OSTLER_DIR/bin
          OSTLER_HOME        -> $HOME (working directory)
          OSTLER_LOGS        -> $LOGS_DIR (~/.ostler/logs by default)
          OSTLER_PYTHON_PATH -> the imessage-source venv python3
          PWG_CONVO_CMD      -> absolute pwg-convo invocation
                                (the CM048 venv python -m src.cli, or
                                the installed console script)

        British English throughout.
    -->
    <key>Label</key>
    <string>com.creativemachines.ostler.imessage-bundle</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>OSTLER_BIN/imessage-bundle-tick.sh</string>
    </array>

    <!--
        Fifteen minutes. iMessage delivers in near-real-time; a
        quarter-hour cadence keeps recent conversations fresh in the
        assistant's memory without hammering chat.db. The watermark
        in ~/.ostler/workspace/imessage_source_state.json means each
        tick only processes sessions with genuinely new messages.
    -->
    <key>StartInterval</key>
    <integer>900</integer>

    <!--
        RunAtLoad fires once at login / post-install. The wrapper
        clamps the read window via --since-days so a fresh install
        does not try to bundle the user's entire iMessage history in
        one tick.
    -->
    <key>RunAtLoad</key>
    <true/>

    <key>WorkingDirectory</key>
    <string>OSTLER_HOME</string>

    <key>StandardOutPath</key>
    <string>OSTLER_LOGS/imessage-bundle.log</string>
    <key>StandardErrorPath</key>
    <string>OSTLER_LOGS/imessage-bundle.err</string>

    <!--
        OSTLER_PYTHON selects the imessage-source venv interpreter.
        PWG_CONVO_CMD tells the dispatcher how to invoke CM048's CLI
        (rendered to the CM048 venv at install time). Without these
        the tick falls through to a stock system python with no
        site-packages -- the same launch-blocker shape the email
        ingest agent guards against.
    -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>OSTLER_PYTHON</key>
        <string>OSTLER_PYTHON_PATH</string>
        <key>PWG_CONVO_CMD</key>
        <string>PWG_CONVO_CMD_VALUE</string>
        <!--
            Rendered by the installer (_install_conversation_feed) so the
            wrapper ships verbatim: SOURCE_DIR is the parent of the
            services/ namespace the wrapper cd's into; USER_DISPLAY_NAME
            labels the operator's own messages as "You".
        -->
        <key>OSTLER_SOURCE_DIR</key>
        <string>OSTLER_SOURCE_DIR_VALUE</string>
        <key>OSTLER_USER_DISPLAY_NAME</key>
        <string>OSTLER_USER_DISPLAY_NAME_VALUE</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>Nice</key>
    <integer>5</integer>
</dict>
</plist>
